# deploy/rcape-bot.service
# RC Ape @-mention bot — systemd unit (reboot-survivable).
# Install (on Malone):
#   sudo cp deploy/rcape-bot.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now rcape-bot
# Logs: journalctl -u rcape-bot -f
#
# Adjust User/WorkingDirectory to match the checkout. The bot reads creds from
# the project .env (loaded by the `bot` npm script via node --env-file).

[Unit]
Description=RC Ape @-mention bot (provisions court-docket repos on request)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=rcape
WorkingDirectory=/opt/rcape
ExecStart=/usr/bin/npm run bot
Restart=always
RestartSec=10
# The poll loop is resilient to transient errors; restart covers crashes/reboots.

# Sandbox hardening. The bot only needs to write its durable stores under data/.
# ProtectSystem=strict mounts the whole filesystem read-only except the paths in
# ReadWritePaths, so a compromised process can't tamper with the host or the
# checkout's code. PrivateTmp gives an isolated /tmp; NoNewPrivileges blocks
# setuid/capability escalation. npm's cache/logs are redirected into the
# PrivateTmp so they don't need write access to $HOME (which strict makes RO).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/rcape/data
Environment=npm_config_cache=/tmp/.npm npm_config_update_notifier=false

[Install]
WantedBy=multi-user.target
